# usr_lib_systemd_system_encheres-publiques.timer
[Unit]
Description=Enchères Publiques

[Timer]
OnCalendar=00/2:30

[Install]
WantedBy=timers.target