# archlinux-package/usr_lib_systemd_system_encheres-publiques.service
[Unit]
Description=Enchères Publiques

[Service]
Type=oneshot
User=encheres-publiques
WorkingDirectory=~
ExecStart=/usr/bin/encheres-publiques